# templates/default/update-motd.d/35-middlewares.erb
#!/bin/sh
nginx_ver=`rpm -q --qf '%{version}' nginx`
<% if @hhvm %>
hhvm_ver=`rpm -q --qf '%{version}' hhvm`
<% elsif @mod_php7 %>
mod_php_ver=`rpm -q --qf '%{version}' php70-php`
<% else %>
php_ver=`rpm -q --qf '%{version}' php`
<% end %>
mysql_ver=`rpm -q --qf '%{version}' Percona-Server-server-56`

cat << EOF

<% if @hhvm %>
 Nginx ${nginx_ver}, HHVM ${hhvm_ver}, Percona MySQL ${mysql_ver}, $(/usr/local/bin/wp --allow-root --version)
<% elsif @mod_php7 %>
 Nginx ${nginx_ver}, MOD_PHP ${mod_php_ver}, Percona MySQL ${mysql_ver}, $(/usr/local/bin/wp --allow-root --version)
<% else %>
 Nginx ${nginx_ver}, PHP ${php_ver}, Percona MySQL ${mysql_ver}, $(/usr/local/bin/wp --allow-root --version)
<% end %>

 amimoto     http://www.amimoto-ami.com/
 digitalcube https://en.digitalcube.jp/

EOF

